# templates/init/statsite.conf
# statsite
#

description     "statsite StatsD service"

start on virtual-filesystems
stop on runlevel [06]

respawn
respawn limit 5 30
limit nofile 65550 65550

console log

script
  exec /usr/local/sbin/statsite -f /etc/statsite.conf
end script
