# templates/systemd/graphite-build-index.service
[Unit]
Description=Run graphite-build-index

[Service]
Type=simple
ExecStart=/usr/local/bin/carbon-build-index

[Install]
WantedBy=multi-user.target
