# templates/systemd/graphite-build-index.timer
[Unit]
Description=Run graphite-build-index hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
